# config/board/qemu/versatilepb/rootfs_overlay/etc/profile.d/exportVariables.sh
#!/bin/sh
#
# Export Variables
#

export QWS_KEYBOARD=linuxinput:/dev/input/event0
export PS1='$(whoami)@$(hostname):\w$ '
